# System/dev/services-vm/beanstalkd-admin/beanstalk-aurora.conf
[program:beanstalk-aurora]
process_name=%(program_name)s_%(process_num)02d
command=/opt/beanstalk-aurora/aurora
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/aurora-beanstalk.log
